# debian/acestep.cpp-sycl.service
[Unit]
Description=ACE-Step music generation server (SYCL)
After=network.target

[Service]
Type=simple
User=acestepcpp
Group=acestepcpp
ExecStart=/usr/bin/ace-server --host 0.0.0.0 --port 8085 --models /srv/acestep.cpp/models --adapters /srv/acestep.cpp/adapters
Environment=GGML_BACKEND=SYCL0
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/acestep.cpp
PrivateTmp=true

[Install]
WantedBy=multi-user.target
